# 1-usb.rules
ACTION=="add", SUBSYSTEM=="block", \
TAG+="systemd", \
RUN+="/bin/systemctl --no-block restart usb_sync.service"
